# script/crontab/crontab.sh
#!/usr/bin/env bash
cat crontab.txt >> /var/spool/cron/root